# scripts/os/desktop-kubernetes.service
[Unit]
Description=Desktop Kubernetes
After=default.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/bash /root/desktop-kubernetes.sh
TimeoutStartSec=0

[Install]
WantedBy=default.target
